# supervisord.conf
[supervisord]
nodaemon=true
user=lumitrade
logfile=/var/log/lumitrade/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info

[program:lumitrade-engine]
command=python -m lumitrade.main
directory=/app
user=lumitrade
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=10
stdout_logfile=/var/log/lumitrade/engine-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/lumitrade/engine-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1"

[program:health-server]
command=python -m lumitrade.infrastructure.health_server
directory=/app
user=lumitrade
autostart=true
autorestart=true
startretries=3
startsecs=3
stdout_logfile=/var/log/lumitrade/health-stdout.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/var/log/lumitrade/health-stderr.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2
environment=PYTHONUNBUFFERED="1"
