# backend/supervisord.conf
[supervisord]
nodaemon=true
user=lumitrade
logfile=/var/log/lumitrade/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:lumitrade-engine]
command=python -m lumitrade.main
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
environment=PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1",NUMBA_DISABLE_JIT="1"

[program:lumitrade-health]
command=python -m lumitrade.infrastructure.health_server
directory=/app
autostart=true
autorestart=true
startretries=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
