# backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:lumitrade-engine]
command=python -m lumitrade.main
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1",NUMBA_DISABLE_JIT="1"

[program:health-server]
command=python -m lumitrade.infrastructure.health_server
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"
